# chef-sync/files/chef-sync.timer
[Unit]
Description=Execute's Edi's runlist

[Timer]
OnCalendar=*-*-* 00/2:30:00
Persistent=true

[Install]
WantedBy=timers.target
