# cookbooks/chef-repo-git-sync/files/chef-repo-git-sync.timer
[Unit]
Description=Resyncs /root/chef/ to master

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
