# cookbooks/weekly-reboot/files/weekly-reboot.timer
[Unit]
Description=Reboot every Monday at 0300

[Timer]
OnCalendar=weekly
 
[Install]
WantedBy=timers.target
